# detections/sentinel/persistence/suspicious-scheduled-task-creation.yml
title: Suspicious Scheduled Task Creation
id: SENT-PERS-0002
status: experimental
description: Detects scheduled task creation via schtasks or Task Scheduler PowerShell that may establish persistence.
author: Adam Ring
date: 2026-03-06
modified: 2026-03-06
platform: microsoft_sentinel
query_language: kql
logsource:
  product: windows
  category: process_creation
query: |
  DeviceProcessEvents
  | where FileName in~ ("schtasks.exe","powershell.exe")
  | where ProcessCommandLine has_any ("/create","Register-ScheduledTask","New-ScheduledTask")
  | project Timestamp, DeviceName, AccountName, FileName, ProcessCommandLine, InitiatingProcessFileName
severity: high
risk_score: 72
tactics:
  - Persistence
techniques:
  - T1053.005
kill_chain_phases:
  - Persistence
data_sources:
  - DeviceProcessEvents
falsepositives:
  - Legitimate enterprise task deployment
  - Admin automation and maintenance jobs
triage:
  - Review the task name, trigger, and action
  - Determine whether the creator is expected to manage tasks
  - Look for follow-on persistence or lateral movement activity
validation:
  - Create a benign scheduled task via schtasks and PowerShell
lifecycle: experimental
owner: Detection Engineering
tags:
  - attack.persistence
  - attack.t1053.005
  - sentinel
  - persistence
